# templates/god.conf.erb
description "God process monitoring"

# When to start and stop  the service
start on runlevel [2345]
stop on runlevel [!2345]

# Automatically restart
respawn

# Tell upstart that the process will detach itself to the background
expect fork

# Start god
exec <% if @use_rvm == 'true' %>/usr/local/rvm/bin/ruby_god<% else %>/usr/bin/god<% end %> -l /var/log/god.log -c <%= @god_master_conf %> -P /var/run/god.pid
