# templates/god.conf.erb
description "God process monitoring"

# When to start and stop  the service
start on runlevel [2345]
stop on runlevel [016]

# Automatically restart
respawn

# Tell upstart that the process will detach itself to the background
expect fork

# Start god
exec <% if @use_rvm == 'true' %><%= @rvm_gem_path %>/<%= @ruby_version %>@<%= @rvm_gemset %>/bin/god<% else %>/usr/bin/god<% end %> -l /var/log/god.log -c <%= @god_master_conf %>
